# etc/sboui-sync.cron.in
#!/bin/sh
#
# This script automatically syncs the SBo local repo using cron

@CMAKE_INSTALL_PREFIX@/sbin/sboui -s
